# src/static/static/home/yi-hack-v4/etc/crontabs/root
12 2 * * * /home/yi-hack-v4/script/check_update.sh
